# proxy-tun.service
[Unit]
Description=SOCKS5 TUN Proxy
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/proxy-tun.sh
ExecStop=/usr/local/bin/proxy-tun.sh stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
